# alpine-consul/root/etc/services.d/consul/run
#!/usr/bin/with-contenv sh

BIND=`getent hosts ${HOSTNAME} | awk '{print $1}'`

exec s6-setuidgid consul consul agent -config-dir /etc/consul.d/bootstrap -bind ${BIND};
