# alpine-consul-python-example/etc/services.d/app/run
#!/usr/bin/env sh

cd /app

exec python server.py