# jupyter-tensorflow/etc/services.d/app/run
#!/usr/bin/env sh

exec jupyter notebook /notebooks --port=8888 --ip=0.0.0.0 --no-browser